# scripts/launchd/com.joplin-brain.chroma.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.joplin-brain.chroma</string>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/scripts/launchd/shims/joplin-llm-wiki-chroma-server</string>
    <string>__REPO_ROOT__/scripts/launchd/run-chroma.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/joplin-llm-wiki/chroma.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/joplin-llm-wiki/chroma.err.log</string>
</dict>
</plist>
